# web2py/web2py/lighttpd-sv.conf
[program:lighttpd]
directory=/opt/web2py
command=/usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=-

